# udev-mic.rules
# do not edit this file, it will be overwritten on update
# initramfs:default

# MIC SCIF
KERNEL=="scif", ACTION=="add", NAME="mic/%k",MODE="0666", RUN+="/bin/chmod og+x /dev/mic"
KERNEL=="ctrl", ACTION=="add", NAME="mic/%k", MODE="0666"
